# pihole/s6/debian-root/etc/cron.d/gravity-sync
@hourly root /bin/bash /usr/local/bin/gravity-sync